# data/etc/init.d/ostree-finalize-staged
#!/sbin/openrc-run

description="OSTree Finalize Staged Deployment"

depend() {
    need localmount
    before killprocs
}

start() {
    ebegin "OSTree finalize staged deployment"
    if ! [ -f /run/ostree-booted ]; then
        eend 0 "System is not booted from an OSTree deployment"
        return 0
    fi

    NEW_SYSLINUX_CFG="/boot/loader/syslinux.cfg"
    if ! [ -f ${NEW_SYSLINUX_CFG} ]; then
        eend 0
        return 0
    fi

    /usr/bin/ostree admin finalize-staged
    mv ${NEW_SYSLINUX_CFG} /boot/syslinux/
    eend $?
}
